# resources/sysinit/gentoo/init.d/syslogng-web
#!/sbin/runscript
# $Id$

depend() {
	need net
	use logger
}

start() {
	
	checkpath -d -m 0750 ${SYSLOGNGWEB_RUN}
		
	ebegin "Starting ${SVCNAME}"
	
	PORT=${SYSLOGNGWEB_PORT:-3000} NODE_ENV=${SYSLOGNGWEB_ENV:-production} start-stop-daemon \
		--background --start --make-pidfile \
		--pidfile ${SYSLOGNGWEB_RUN:-/var/run}/${SVCNAME}.pid \
		--stdout ${SYSLOGNGWEB_LOG:-/var/log/syslogng-web.log} \
		--exec node ${SYSLOGNGWEB_DIR:-/opt/syslogng-web}/app.js
	
	eend $?

}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --pidfile ${SYSLOGNGWEB_RUN:-/var/run}/${SVCNAME}.pid
	eend $?
}
