# resources/sysinit/gentoo/init.d/syslogng-web
#!/sbin/runscript

depend() {
	need net
	use logger
}

start() {
	
	checkpath -d -m 0750 ${SYSLOGNGWEB_RUN}
		
	ebegin "Starting ${SVCNAME}"
	PORT=${SYSLOGNGWEB_PORT} start-stop-daemon --background --start --make-pidfile \
	--pidfile ${SYSLOGNGWEB_RUN:-/var/run/syslogng-web}/${SVCNAME}.pid \
	--exec $node ${SYSLOGNGWEB_DIR:-/usr/local/share/syslogng-web}/app.js
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --pidfile ${SYSLOGNGWEB_RUN:-/var/run/syslogng-web}/${SVCNAME}.pid
}
